# reporting/hier/etc/cron.d/untangle-reports-nightly
0 2 * * * root /usr/share/untangle/bin/reporting-generate-reports.py -d $(date "+%Y-%m-%d") > /dev/null 2>&1
