# uvm/hier/etc/cron.d/untangle-reports-incremental
*/5 */1 * * * root /usr/share/untangle/bin/reporting-generate-reports.py -m -i -b 1 -r 1 > /dev/null 2>&1
