# mvvm/hier/etc/cron.daily/mvvm-cron
#!/bin/sh

/usr/bin/update-reports &> /dev/null
