# debian/mvvm.files
etc/cron.daily/mvvm-cron
etc/default/mvvm
etc/init.d/mvvm
usr/bin/*
usr/lib/mvvm/libalpine.so
usr/share/java/*
usr/share/java/mvvm/*
usr/share/metavize/conf/*
usr/share/metavize/conf/etc
usr/share/metavize/conf/etc/network
usr/share/metavize/lib/*
usr/share/metavize/log
usr/share/metavize/mvvm_restart.sh
usr/share/metavize/schema/mvvm/*
usr/share/metavize/ssh_disable.sh
usr/share/metavize/ssh_enable.sh
usr/share/metavize/web/*
usr/share/metavize/networking/*
usr/share/metavize/tmp/*

