# debian/untangle-vm.postrm
#! /bin/bash

if [ "$1" = "purge" ] ; then
    rm -f /etc/init.d/untangle-vm
    systemctl disable untangle-vm > /dev/null
fi

exit 0
