# debian/untangle-vm.postrm
#!/bin/sh

if [ "$1" = "purge" ] ; then
    rm -f /etc/init.d/untangle-vm
    update-rc.d untangle-vm remove > /dev/null
fi

