# debian/untangle-vm.postrm
#!/bin/sh

if [ "$1" = "purge" ] ; then
    rm -f /etc/init.d/mvvm
    update-rc.d mvvm remove > /dev/null
fi

