# uvm/hier/etc/cron.daily/untangle-vm-cron
#! /bin/bash

exec > /dev/null 2>&1

# remove old temp files older than 7 days (safety mechanism)
/usr/bin/find /tmp -mtime +7 -name "*.tmp" | /usr/bin/xargs -r /bin/rm

exit 0
